# supervisord.ini
[program:drf-example]
command = uwsgi --ini uwsgi.ini
directory=/srv/drf-example
stdout_logfile=/srv/drf-example/logs/stdout.log
autostart=true
autorestart=true
user=root
startsecs=3