# supervisord.ini
[program:drf-example]
command = uwsgi --ini uwsgi.ini
directory=/srv/drf-example
stdout_logfile=/srv/drf-example/logs/stdout.log
autostart=true
autorestart=true
user=root
startsecs=3
priority=998

[program:drf-example-celery]
command = celery -A demo worker -l INFO
directory=/srv/drf-example
stdout_logfile=/srv/drf-example/logs/celery.log
autostart=true
autorestart=true
user=root
startsecs=3
priority=999

[program:drf-example-flower]
command = celery -A demo flower --port=30555
directory=/srv/drf-example
stdout_logfile=/srv/drf-example/logs/flower.log
autostart=true
autorestart=true
user=root
startsecs=3
priority=997