# supervisord.ini
[program:edrf]
command = uwsgi --ini uwsgi.ini
directory=/srv/edrf
stdout_logfile=/srv/edrf/logs/stdout.log
autostart=true
autorestart=true
user=root
startsecs=3
priority=998

[program:edrf-celery]
command = celery -A server worker -l INFO
directory=/srv/edrf
stdout_logfile=/srv/edrf/logs/celery.log
autostart=true
autorestart=true
user=root
startsecs=3
priority=999

[program:drf-flower]
command = celery -A server flower --port=30555
directory=/srv/edrf
stdout_logfile=/srv/edrf/logs/flower.log
autostart=true
autorestart=true
user=root
startsecs=3
priority=997
